# robot_control/auto_start_file/after-resume.target
[Unit]
Description=Custom target triggered after resume

