# robot_control/auto_start_file/screen_unlock.target
[Unit]
Description=Custom target triggered after resume 2
DefaultDependencies=no
StopWhenUnneeded=yes
